# qpay-client_0.2-1/etc/systemd/system/qpay-http-bridge.service
[Unit]
Description=Qpay HTTP bridge

[Service]
Type=simple
ExecStart=/usr/bin/qpay -b

[Install]
WantedBy=multi-user.target
